# linux/wikiscraper.service
[Unit]
Description=WikiScraper service
[Service]
ExecStart=/home/discsel/bin/wikiscraper.sh
User=discsel
Restart=always
[Install]