# linux/wikiscraper.service
[Unit]
Description=WikiScraper service
[Service]
ExecStart=/home/leizh/bin/wikiscraper.sh
User=leizh
Restart=always
[Install]